# autostart/rescuetime.desktop
[Desktop Entry]
Name=Rescuetime applet
TryExec=rescuetime
Exec=rescuetime &
Terminal=false
Type=Application
StartupNotify=true
OnlyShowIn=awesome
